# deploy/systemd/poverlay-web.service
[Unit]
Description=POVerlay Next.js web service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/poverlay
EnvironmentFile=-/etc/poverlay/poverlay.env
Environment=NODE_ENV=production
Environment=HOME=/home/ubuntu
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/poverlay/apps/web/.next/standalone/apps/web/server.js
Restart=always
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
